# assets/supervisord.conf
[supervisord]
nodaemon=true
user=root

;logfile=/var/log/supervisor/supervisord.log  ; (main log file;default $CWD/supervisord.log)
;logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
;logfile_backups=10          ; (num of main logfile rotation backups;default 10)
;loglevel=info               ; (log level;default info; others: debug,warn,trace)
;pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
;nodaemon=false              ; (start in foreground if true;default false)
;minfds=1024                 ; (min. avail startup file descriptors;default 1024)
;minprocs=200                ; (min. avail process descriptors;default 200)

[program:sshd]
command=/usr/sbin/sshd -D
#command=/usr/sbin/sshd -D
#command=service sshd start
#user=root

[program:apache2]
#command=apachectl -D "FOREGROUND" -k start
;redirect_stderr=true
command=/usr/local/apache2/bin/httpd -D FOREGROUND
#user=root
;command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/local/apache2/apache2 -DFOREGROUND"
